# deploy/ticketu-analyser.service
[Unit]
Description=Analyse les tickets PDF Système U
Documentation=https://github.com/nikkocasa/TicketUanalyser
After=network.target docker.service

[Service]
Type=oneshot
User=ticketanalyser
Group=ticketanalyser
WorkingDirectory=/opt/ticketanalyser

# Petite pause pour s'assurer que le fichier est entièrement écrit avant de lire
ExecStartPre=/bin/sleep 2
ExecStart=/opt/ticketanalyser/venv/bin/python /opt/ticketanalyser/ticket_analyser.py \
    --config /opt/ticketanalyser/config.ini

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ticketu-analyser
