# deploy/ticketu-dashboard.service
[Unit]
Description=Dashboard Streamlit TicketU
Documentation=https://github.com/nikkocasa/TicketUanalyser
After=network.target docker.service

[Service]
Type=simple
User=ticketanalyser
Group=ticketanalyser
WorkingDirectory=/opt/ticketanalyser
ExecStart=/opt/ticketanalyser/venv/bin/streamlit run dashboard.py \
    --server.headless=true \
    --server.port=8501 \
    --server.address=0.0.0.0
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ticketu-dashboard

[Install]
WantedBy=multi-user.target
